# udev/99-webcam-settings.rules
SUBSYSTEM=="video4linux", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="082d", RUN+="/usr/bin/v4l2-ctl -d /dev/%k --set-ctrl=power_line_frequency=1"
